# itamae/files/etc/profile.d/caffe.sh
caffe_prefix=/usr/local/caffe/distribute

export PATH=$caffe_prefix/bin:$PATH
export LD_LIBRARY_PATH=$caffe_prefix/lib64:$LD_LIBRARY_PATH
export PYTHONPATH=$caffe_prefix/python:$PYTHONPATH
