# itamae/files/etc/profile.d/cuda.sh
export cuda_prefix=/usr/local/cuda-6.5

export PATH=$cuda_prefix/bin:$PATH
export LD_LIBRARY_PATH=$cuda_prefix/lib64:$LD_LIBRARY_PATH
