# itamae/files/etc/profile.d/torch.sh
torch_prefix=/usr/local/torch

export PATH=$torch_prefix/bin:$PATH
export LD_LIBRARY_PATH=$torch_prefix/lib:$LD_LIBRARY_PATH
